# home/.config/systemd/user/clipboard.socket
[Unit]
Description=Clipboard Socket
PartOf=clipboard.service

[Socket]
ListenStream=127.0.0.1:9998

[Install]
WantedBy=sockets.target
